# deployment/chatwoot-worker.target
[Unit]
Description=nerix Background Worker
Wants=nerix-worker.1.service

[Install]
WantedBy=multi-user.target