# deployment/chatwoot.target
[Unit]
Wants=nerix-web.1.service nerix-worker.1.service

[Install]
WantedBy=multi-user.target