# deployment/chatwoot-web.target
[Unit]
Description=clord Web Server
Wants=clord-web.1.service

[Install]
WantedBy=multi-user.target