# deployment/chatwoot.target
[Unit]
Wants=Clord-web.1.service Clord-worker.1.service

[Install]
WantedBy=multi-user.target