# deployment/chatwoot-worker.1.service
[Unit]
Requires=network.target
PartOf=nerix.target

[Service]
Type=simple
User=nerix
WorkingDirectory=/home/nerix/nerix

ExecStart=/bin/bash -lc 'dotenv bundle exec sidekiq -C config/sidekiq.yml'

Restart=always
RestartSec=1
TimeoutStopSec=30
KillMode=mixed
StandardInput=null
SyslogIdentifier=%p
LimitNOFILE=65536

MemoryMax=1.2G
MemoryHigh=infinity
MemorySwapMax=0
OOMPolicy=stop

Environment="PATH=/home/nerix/.rvm/gems/ruby-3.4.4/bin:/home/nerix/.rvm/gems/ruby-3.4.4@global/bin:/home/nerix/.rvm/rubies/ruby-3.4.4/bin:/home/nerix/.rvm/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin:/home/nerix/.rvm/bin:/home/nerix/.rvm/bin"
Environment="PORT=3000"
Environment="RAILS_ENV=production"
Environment="NODE_ENV=production"
Environment="RAILS_LOG_TO_STDOUT=true"
Environment="GEM_HOME=/home/nerix/.rvm/gems/ruby-3.4.4"
Environment="GEM_PATH=/home/nerix/.rvm/gems/ruby-3.4.4:/home/nerix/.rvm/gems/ruby-3.4.4@global"
