# deployment/chatwoot.target
[Unit]
Wants=clord-web.1.service clord-worker.1.service

[Install]
WantedBy=multi-user.target