# deployment/chatwoot-worker.target
[Unit]
Description=clord Background Worker
Wants=clord-worker.1.service

[Install]
WantedBy=multi-user.target